# Codes/server/datastore/build/omp/monitor-datastore.service
[Unit]  
Description=monitor-datastore
After=network.target  
   
[Service]  
Type=oneshot
RemainAfterExit=yes  
ExecStart=/usr/bin/su -l root /opt/uyun/monitor/datastore/bin/bat-datastore.sh start
ExecStop=/usr/bin/su -l root /opt/uyun/monitor/datastore/bin/bat-datastore.sh stop
   
[Install]  
WantedBy=multi-user.target
